# fs-nginx/etc/services/nginx/run
#!/usr/bin/env sh

# Wait up to 10 seconds for php-fpm to start.
sv -w 10 check /etc/services/php-fpm

exec nginx
