# base/services/fpm/run
#!/command/execlineb -P

with-contenv php-fpm