# base/services/set-user-id/run
#!/usr/bin/env sh

set -e

usermod -u "$PUID" app
groupmod -g "$PGID" app