# base/services/set-user-id/up
with-contenv sh /etc/s6-overlay/s6-rc.d/set-user-id/run